# zcm-nginx-ingress-controller/httpfiles/gateway-ini/zcm-nginx-gateway.ini
[program:zcm-nginx-gateway]
command=/bin/bash -c \
        '. /zcm/supervisor/zcm-env.profile; \
        [ `docker ps -a -f "name=zcm-nginx-gateway"| wc -l` -gt 1 ] && \
        exec docker start -a zcm-nginx-gateway || \
        exec docker run --name zcm-nginx-gateway \
	--restart=unless-stopped \
        --privileged --network host \
	-m 4096m --memory-swap=4096m \
        -v /etc/localtime:/etc/localtime:ro \
        -v /etc/timezone:/etc/timezone:ro \
        -v /zcm/zcm-gateway/application.properties:/root/conf/application.properties \
        -v /tmp/zlogs:/tmp/zlogs \
        -v /zcm/zcm-gateway/data:/root/data \
        -e ZCM_GATEWAY_PORT=${ZCM_GATEWAY_PORT} \
	-e ZCM_GATEWAY_SSL_PORT=${ZCM_GATEWAY_SSL_PORT} \
        -e ZCM_SERVER_VIP=${ZCM_SERVER_VIP} \
        -e ZCM_SERVER_PORT=${ZCM_SERVER_PORT} \
        -e ZCM_SERVER_SSL_PORT=${ZCM_SERVER_SSL_PORT} \
        -e GATEWAY_SERVER_IP=${K8S_MASTER_URL} \
	-e APM_INIT=${ZCM_APM_INIT}  \
	-e CLOUD_APP_NAME="paas_zcm-nginx-gateway"  \
        -d \
        ${ZCM_HARBOR}/${ZCM_HARBOR_PROJECT}/zcm-gateway:${ZCM_GATEWAY_TAG}; \
        docker start -a zcm-nginx-gateway'
process_name= zcm-nginx-gateway
numprocs=1
autostart=true
autorestart=true
startretries = 6
user=root
stdout_logfile=/var/log/zcm/zcm-nginx-gateway.log
stdout_logfile_maxbytes=10MB
stderr_logfile=/var/log/zcm/zcm-nginx-gateway.log
stderr_logfile_maxbytes=10MB
stopsignal=INT
stopwaitsecs=20
killasgroup=true
stopasgroup=true